# acestream.2/dockerfile-ffmpeg/ffmpeg.conf
[program:ffmpeg]
command=/acelink/hls.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/ffmpeg.err.log
stdout_logfile=/var/log/ffmpeg.out.log
